# dockter-tom/supervisord.conf
[supervisord]
logfile                     =/var/log/supervisor/supervisord.log    ; (main log file;default $CWD/supervisord.log)
childlogdir                 = /var/log/supervisor                   ; (log directory for child process logs)
logfile_maxbytes            =50MB                                   ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups             =2                                      ; (num of main logfile rotation backups;default 10)
loglevel                    =debug                                  ; (log level;default info; others: debug,warn,trace)
pidfile                     =/var/run/supervisord.pid               ; (supervisord pidfile;default supervisord.pid)
nodaemon                    =true                                   ; (start in foreground if true;default false)

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

;; nmap
[program:nmap]
command ="/usr/bin/nmap -p1-1024 localhost  >>/opt/results 2>&1" ;; $1
autorestart=false

;; nikto
[program:nikto]
command = "perl /opt/nikto/program/nikto.pl -h localhost -p 80,88,443 >>/opt/results 2>&1" ;; $1
autorestart=false

;; sqlmap
[program:sqlmap]
command = "python /opt/sqlmap/sqlmap.py -a localhost >>/opt/results 2>&1" ;; $1
autorestart=false

;; gauntlt
;; [program:gauntlt]
;; command = /usr/local/bin/gaultlt /opt/attack/*attacks >>/opt/results 2>&1 #$1

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket