# supervisord.conf
[supervisord]
logfile                     =/var/log/supervisor/supervisord.log    ; (main log file;default $CWD/supervisord.log)
childlogdir                 = /var/log/supervisor                   ; (log directory for child process logs)
logfile_maxbytes            =50MB                                   ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups             =2                                      ; (num of main logfile rotation backups;default 10)
loglevel                    =debug                                  ; (log level;default info; others: debug,warn,trace)
pidfile                     =/var/run/supervisord.pid               ; (supervisord pidfile;default supervisord.pid)
nodaemon                    =true                                   ; (start in foreground if true;default false)

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

;; nmap
[program:nmap]
command=/usr/local/bin/nmap -p80, 443, 3306 localhost >> /opt/results

;; nikto
[program:nikto]
command=perl nikto -h localhost -p 80,88,443

;; sqlmap
[program:httpd]
command=/opt/sqlmap/sqlmap.py localhost

;; sshd
[program:sshd]
command=/usr/sbin/sshd -D 


; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket